# deploy/host/proxsync-agent.service
[Unit]
Description=ProxSync Backup Agent
Documentation=https://github.com/stayfr0sty141/proxsync/blob/main/docs/ARCHITECTURE.md
After=network-online.target pve-cluster.service
Wants=network-online.target
ConditionPathExists=/usr/bin/vzdump

[Service]
Type=exec
ExecStart=/opt/proxsync-agent/venv/bin/python -m app
WorkingDirectory=/opt/proxsync-agent
EnvironmentFile=/etc/proxsync-agent/agent.env

# --- Privilege -------------------------------------------------------------
# vzdump, qmrestore and pct restore require root and take dynamic arguments, so a sudoers
# rule permissive enough to be useful would also be permissive enough to abuse. The request
# validator in app/validators is the real security boundary; this unit narrows what a
# compromised agent could reach beyond its intended job.
User=root
Group=root

# --- Sandbox ---------------------------------------------------------------
NoNewPrivileges=yes
ProtectSystem=true
# NOT "full" or "strict": qm and pct write guest configuration to /etc/pve, so a read-only
# /etc breaks every restore. Do not tighten this without testing a restore.
ProtectHome=yes
# Remove ProtectHome if your backup storage is mounted under /home.
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes
RestrictSUIDSGID=yes
RestrictNamespaces=no
# NOT restricted: pct restore creates mount and user namespaces for containers.
RestrictRealtime=yes
LockPersonality=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
# AF_NETLINK is required: LVM and udev use it during snapshot backups.
PrivateTmp=no
# NOT private: vzdump stages archives in the configured tmpdir and expects it to be shared.
PrivateDevices=no
# NOT private: LVM/ZFS snapshot backups need access to block devices.
MemoryDenyWriteExecute=no
# NOT denied: helper processes spawned by vzdump may JIT.

# --- Network reachability --------------------------------------------------
# Note on systemd IP filtering: IPAddressDeny=any blocks ALL inbound & outbound connections.
# Because the agent invokes rclone for Google Drive syncs (which connect to dynamic outbound
# Google endpoints), systemd IP filtering is not used here.
# Defence in depth is enforced via:
# 1) Bind IP (PROXSYNC_AGENT_BIND_HOST)
# 2) Host nftables / iptables rules
# 3) Mutual TLS (mTLS) with dedicated CA
# 4) HMAC request signature verification
# 5) App layer PROXSYNC_AGENT_ALLOWED_CLIENT_NETWORKS filtering

# --- Managed directories ---------------------------------------------------
StateDirectory=proxsync-agent
StateDirectoryMode=0750
LogsDirectory=proxsync-agent
LogsDirectoryMode=0750
ConfigurationDirectory=proxsync-agent
ConfigurationDirectoryMode=0750

# --- Resources and lifecycle ----------------------------------------------
# The agent process and its spawned children (vzdump, qmrestore, pct, rclone, sha256sum)
# share this systemd cgroup. Memory limits are set to allow heavy backup/sync operations.
MemoryHigh=1G
MemoryMax=2G
TasksMax=256
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
# Default KillMode=control-group: a stop kills vzdump too. The agent marks such tasks
# "interrupted" on its next start rather than guessing they succeeded.

StandardOutput=journal
StandardError=journal
SyslogIdentifier=proxsync-agent

[Install]
WantedBy=multi-user.target
